# davfs2/mnt-nextcloud.automount
[Unit]
Description=Mount WebDAV Service

[Automount]
Where=/mnt/nextcloud
TimeoutIdleSec=300

[Install]
WantedBy=remote-fs.target
