# davfs2/mnt-nextcloud.automount
[Unit]
Description=Mount WebDAV Service

[Automount]
Where=/mnt/nextcloud
TimeoutIdleSec=0

[Install]
WantedBy=remote-fs.target
